# external/systemd/system/goldstone-server.service
[Unit]
Description=Goldstone Server
Requires=docker.service
After=docker.service

[Service]
Restart=always
User=root
Group=root
ExecStart=/usr/bin/docker-compose -f /opt/golstone/docker/docker-compose.yml up
ExecStop=/usr/bin/docker-compose -f /opt/goldstone/docker/docker-compose.yml stop

[Install]
WantedBy=local.target
